# run.crontab.sample
# crontab sample file to schedule the wearther update of the ThingSpeak server 4 times per day

0 9,14,19,23 * * * /usr/bin/python2.7 /home/swampi/SwampertReloaded/telegramBot.py
